# stable-diffusion.service
[Unit]
Description=Stable-Diffusion
After=network-online.target

[Service]
Type=simple

User=username
Group=username
ExecStart=/opt/diffusion/stable-diffusion-webui/start.sh
WorkingDirectory=/opt/diffusion/
Restart=on-failure

[Install]
WantedBy=multi-user.target